# dnsupdate.timer
[Unit]
Description=Dynamic DNS update timer

[Timer]
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
